# debian/vyatta-webgui.postrm
#!/bin/bash

update-rc.d -f oa-autostart remove

