# debian/vyatta-webgui.postrm
#!/bin/bash

update-rc.d oa-autostart remove

